# supervisor.conf
; FreeIOT supervisor config file.
;
; For FreeIOT Program, Author Noah Gao
; Updated at: 2018-02-24

[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file

[supervisord]
logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=false               ; start in foreground if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:freeiot-sample]
command=/home/noah/freeiot/venv/bin/python3 ./manage.py     ; supervisor启动命令
directory=/home/noah/freeiot                                ; 项目的文件夹路径
startsecs=0                                                 ; 启动时间
stopwaitsecs=0                                              ; 终止等待时间
autostart=true                                              ; 是否自动启动
autorestart=true                                            ; 是否自动重启
stdout_logfile=/home/noah/server/log/freeiot-sample.log     ; log 日志
stderr_logfile=/home/noah/server/log/freeiot-sample.err     ; error 日志
